# assets/supervisor/mocks.ini
[program:mocks]
command=python /usr/bin/mocks.py
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
